# deploy/irina-crm-backup.timer
[Unit]
Description=Daily Irina CRM backup

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=irina-crm-backup.service

[Install]
WantedBy=timers.target
